# provisioning/roles/ipython_virtualenv/files/ipython_server.ini
[program:ipython_server]
directory=/home/ipython
user=ipython
command=/home/ipython/env/bin/ipython notebook --no-browser --profile=nbserver
stderr_logfile = /tmp/ipython-notbook-server-stderr.log
stdout_logfile = /tmp/ipython-notbook-server-stdout.log
